# etc/systemd/user/brew-update.timer
[Unit]
Description=Run Homebrew Auto Update daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
